# docker/s6-rc.d/browserbridge/run
#!/command/with-contenv sh
# s6 service: BrowserBridge (combined WebSocket + MCP server)
# The Node.js entry point handles token auto-generation and startup banner.

cd /app
exec node dist/bin/browserbridge.js